# Controller/systemctl/motor.service
[Unit]
Description=Start Motor

[Service]
ExecStart=/usr/bin/python3 /home/pi/proj/startMotor.py

[Install]
WantedBy=multi-user.target